# snat.service
[Unit]
Description = SNAT via ENI ens6

[Service]
ExecStart = /opt/nat/snat.sh
Type = oneshot

[Install]
WantedBy = multi-user.target
